# deploy/prefect.service
[Unit]
Description=Prefect Server
After=network.target

[Service]
Type=simple
User=crypto-research
Group=crypto-research
WorkingDirectory=/opt/crypto-research
Environment="PATH=/opt/crypto-research/.venv/bin:/usr/local/bin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PREFECT_SERVER_API_HOST=0.0.0.0"
Environment="PREFECT_API_URL=http://localhost:4200/api"

ExecStart=/opt/crypto-research/.venv/bin/prefect server start --host 0.0.0.0
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
MemoryLimit=2G
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prefect-server

[Install]
WantedBy=multi-user.target
